# minimax-collector.service
[Unit]
Description=MiniMax Coding Plan Collector
After=network.target

[Service]
Type=simple
User=evi1
WorkingDirectory=/home/evi1/minimax-coding-plan-check
ExecStart=/usr/bin/python3 /home/evi1/minimax-coding-plan-check/collector.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
